# templates/navidrome.service.epp
<%- | Optional[String]  $systemd_description,
      Optional[String]  $systemd_path,
      Optional[String]  $systemd_after,
      Optional[String]  $systemd_user,
      Optional[String]  $systemd_group,
      Optional[String]  $systemd_wantedby,
      Optional[String]  $systemd_type,
      Optional[String]  $systemd_navidrome_location,
      Optional[String]  $systemd_settings_location,
      Optional[String]  $systemd_working_dir,
      Optional[Integer] $systemd_timeout_stop_sec,
      Optional[String]  $systemd_killmode,
      Optional[String]  $systemd_restart,
| -%>
[Unit]
Description=<%= $systemd_description %>
After=<%= $systemd_after %>
AssertPathExists=<%= $systemd_path %>

[Install]
WantedBy=<%= $systemd_wantedby %>

[Service]
User=<%= $systemd_user %>
Group=<%= $systemd_group %>
Type=<%= $systemd_type %>
ExecStart=<%= $systemd_navidrome_location %> --configfile "<%= $systemd_settings_location %>"
WorkingDirectory=<%= $systemd_working_dir %>
TimeoutStopSec=<%= $systemd_timeout_stop_sec %>
KillMode=<%= $systemd_killmode %>
Restart=<%= $systemd_restart %>

DevicePolicy=closed
NoNewPrivileges=yes
PrivateTmp=yes
PrivateUsers=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallFilter=~@clock @debug @module @mount @obsolete @reboot @setuid @swap
ReadWritePaths=<%= $systemd_working_dir %>
ProtectSystem=full